# config/taskmaster_conf.ini
[program:random101]
command=sh sec_counter.bash
autostart=true
autorestart=true
retry=auto
dir=./
nbretries=3
stdout=./101.stdout
stderr=./101.stderr
env=
exitcodes=
succodes=
stoptime=10
umask=