# config/taskmaster_conf.ini
[program:random69]
command=sh sec_counter.bash
autostart=true
autorestart=true
retry=auto
dir=./
nbretries=3
stdout=./time_program.stdout
stderr=./time_program.stderr
env=
exitcodes=
succodes=
stoptime=
umask=

