# config/taskmaster_conf.ini
[program:random101]
command=sh sec_counter.bash
autostart=0
autorestart=0
retry=auto
dir=./
nbretries=3
stdout=./101.stdout
stderr=./101.stderr
env=
exitcodes=-1
succodes=-1
stoptime=0
umask=773
