# config/taskmaster_conf.ini
[program:random69]
command=sh sec_counter.bash
autostart=true
autorestart=true
retry=auto
dir=./
nbretries=3
stdout=./69.stdout
stderr=./69.stderr
env=
exitcodes=
succodes=
stoptime=10
umask=